# .config/systemd/user/onedrive.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/onedrive --monitor --local-first
